# 60-msikdb.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="1770", ATTR{idProduct}=="ff00", GROUP="input", RUN+="/usr/bin/systemctl --no-block start msi-kbd.service"
